# packaging/cocoon-agent.service
[Unit]
Description=Cocoon agent (vsock command exec)
Documentation=https://github.com/cocoonstack/cocoon-agent
# vsock is host-local; no network or DNS dependency.

# Crash-loop bound: at most 5 restarts inside 60s before systemd holds the
# unit in failed state. Prevents a wedged hypervisor/viosock state from
# pegging journald with restart churn. StartLimit* live under [Unit] since
# systemd v229 — the [Service]-side spelling is the legacy form.
StartLimitBurst=5
StartLimitIntervalSec=60s

[Service]
Type=simple
User=root
Group=root
# No modprobe ExecStartPre: vhost_vsock is the HOST-side module (driver for
# /dev/vhost-vsock that the hypervisor opens) — loading it inside the GUEST
# is a no-op error. The guest needs virtio_transport_common +
# vmw_vsock_virtio_transport, both of which auto-load on virtio-vsock device
# probe in any modern distro kernel. listenVsockWithRetry already rides out
# the bind window if the device shows up late.
ExecStart=/usr/local/bin/cocoon-agent serve
Environment=AGENT_LOG_LEVEL=info
Restart=always
RestartSec=2s
LimitNOFILE=65536
# KillMode=control-group (default) is correct: in-flight kubectl-exec
# children inherit the unit's cgroup and are torn down with the agent
# on systemctl stop.
#
# The agent intentionally runs as root: it executes commands on behalf
# of host-side callers, so privilege drop would defeat its purpose.
# All other process hardening (PrivateTmp, ProtectSystem, ...) is left
# off to preserve full exec semantics for arbitrary client commands.

[Install]
WantedBy=multi-user.target
